# fs/etc/services.d/portainer/run
#!/usr/bin/with-contenv ash

exec runas /portainer $PORTAINER_ARGS