# fs/etc/services.d/portainer/run
#!/usr/bin/with-contenv ash

exec s6-setuidgid "$PUID:$PGID" /portainer $PORTAINER_ARGS