# fs/etc/services.d/portainer/run
#!/usr/bin/with-contenv sh

s6-setuidgid nonroot

exec /portainer $PORTAINER_ARGS